# systemd/wordcracker-chat.service.d/v2-engine.conf
# Drop-in override: route /api/chat through the v2 planner+router engine by default.
#
# The main unit's ExecStart only forwards ASSISTANT_NAME into `docker compose
# exec`, so a plain `Environment=WC_DEFAULT_ENGINE=v2` on the systemd side
# never reaches chat_server.py inside the container. We reset ExecStart and
# re-add it with `-e WC_DEFAULT_ENGINE=v2`.
#
# Install:
#   sudo cp v2-engine.conf /etc/systemd/system/wordcracker-chat.service.d/
#   sudo systemctl daemon-reload
#   sudo systemctl restart wordcracker-chat
#
# Toggle back to v1 by removing the drop-in:
#   sudo rm /etc/systemd/system/wordcracker-chat.service.d/v2-engine.conf
#   sudo systemctl daemon-reload && sudo systemctl restart wordcracker-chat
#
# Per-request override still works via `?engine=v1` query string or
# `X-WC-Engine: v1` header.

[Service]
ExecStart=
ExecStart=/usr/bin/docker compose exec -T -e ASSISTANT_NAME=Словоёб -e WC_DEFAULT_ENGINE=v2 gutenberg-lab python -u /workspace/scripts/chat_server.py --port 8890
