# systemd/wordcracker-admin.service
[Unit]
Description=wordcracker admin upload server (port 8891, inside gutenberg-lab container)
Documentation=file:///home/claude/wordcracker/scripts/admin_server.py
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=claude
Group=claude
WorkingDirectory=/home/claude/wordcracker
ExecStartPre=/usr/bin/docker compose -f docker-compose.yml up -d gutenberg-lab
# SIGTERM doesn't propagate from `docker compose exec` to the python process
# inside the container, so old admin_server.py instances can survive a restart
# and keep port 8891 bound. Reap any orphan before binding the port. `-`
# prefix swallows the exit code if pkill finds nothing.
ExecStartPre=-/usr/bin/docker compose -f docker-compose.yml exec -T gutenberg-lab pkill -9 -f /workspace/scripts/admin_server.py
# Defensive pyc purge (see chat unit for context).
ExecStartPre=-/usr/bin/docker compose -f docker-compose.yml exec -T gutenberg-lab find /workspace/scripts/__pycache__ -name '*.pyc' -delete
ExecStart=/usr/bin/docker compose -f docker-compose.yml exec -T gutenberg-lab python -u /workspace/scripts/admin_server.py --port 8891
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wc-admin
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
