# systemd/wordcracker-chat.service
[Unit]
Description=wordcracker chat agent (port 8890, inside gutenberg-lab container)
Documentation=file:///home/claude/wordcracker/scripts/chat_server.py
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=claude
Group=claude
WorkingDirectory=/home/claude/wordcracker
Environment=ASSISTANT_NAME=Словоёб
# Ensure the container is up before we try to exec into it. Idempotent.
ExecStartPre=/usr/bin/docker compose up -d gutenberg-lab
# SIGTERM doesn't propagate from `docker compose exec` to the python process
# inside the container, so old chat_server.py instances can survive a restart
# and keep port 8890 bound. Reap any orphan before binding the port. `-`
# prefix swallows the exit code if pkill finds nothing.
ExecStartPre=-/usr/bin/docker compose exec -T gutenberg-lab pkill -9 -f /workspace/scripts/chat_server.py
# Purge stale .pyc files before loading rag_tools — fresh scp of .py won't
# beat an old .pyc if cwd mtime gets touched unexpectedly. Defensive only;
# Python normally invalidates pyc by source mtime, but we hit it once during
# the v1.0f deploy where the new find_book code was masked.
ExecStartPre=-/usr/bin/docker compose exec -T gutenberg-lab find /workspace/scripts/__pycache__ -name '*.pyc' -delete
# Run the chat server in the foreground, attached to systemd. SIGTERM from
# `systemctl stop` propagates to the python process via docker exec.
ExecStart=/usr/bin/docker compose exec -T -e ASSISTANT_NAME=Словоёб gutenberg-lab python -u /workspace/scripts/chat_server.py --port 8890
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wc-chat
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
